# install/wx/systemd/wxconsumer.service
[Unit]
Description=wxconsumer
After=network.target
 
[Service]
Type=forking
ExecStart=/etc/init.d/consumer.sh
 
[Install]
WantedBy=multi-user.target
