# install/wx/systemd/wxip.service
[Unit]
Description=wxip
After=network.target
 
[Service]
Type=forking
ExecStart=/etc/init.d/ip.sh
 
[Install]
WantedBy=multi-user.target
