# install/wx/systemd/wxjob.service
[Unit]
Description=wxjob
After=network.target
 
[Service]
Type=forking
ExecStart=/etc/init.d/job.sh
 
[Install]
WantedBy=multi-user.target
